# config/supervisor-fleet.conf
# supervisord fleet service config
# Drop into /etc/supervisor/conf.d/ on the fleet host
# Then: sudo supervisorctl reread && sudo supervisorctl update

[group:fleet]
programs=dashboard,nexus,harbor,guard,queue,steward

[program:dashboard]
command=python3 -m dashboard.server --port 4046
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/dashboard.err.log
stdout_logfile=/var/log/supervisor/dashboard.out.log
environment=HOME="/home/ubuntu",USER="ubuntu"

[program:nexus]
command=python3 -m nexus.federation --port 4047
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/nexus.err.log
stdout_logfile=/var/log/supervisor/nexus.out.log
environment=HOME="/home/ubuntu",USER="ubuntu"

[program:harbor]
command=python3 -m harbor.server --port 4050
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/harbor.err.log
stdout_logfile=/var/log/supervisor/harbor.out.log
environment=HOME="/home/ubuntu",USER="ubuntu"

[program:guard]
command=python3 -m service_guard.server --port 8899
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/guard.err.log
stdout_logfile=/var/log/supervisor/guard.out.log
environment=HOME="/home/ubuntu",USER="ubuntu"

[program:queue]
command=python3 -m task_queue.server --port 8900
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/queue.err.log
stdout_logfile=/var/log/supervisor/queue.out.log
environment=HOME="/home/ubuntu",USER="ubuntu"

[program:steward]
command=python3 -m steward.server --port 8901
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/supervisor/steward.err.log
stdout_logfile=/var/log/supervisor/steward.out.log
environment=HOME="/home/ubuntu",USER="ubuntu"
